# oxsystemsos/airootfs/etc/xdg/autostart/ox-live-installer.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=OXSystemsOS Live Installer
Comment=Launch OXSystemsOS installer on live GNOME session start
Exec=/usr/local/bin/ox-live-installer --autostart
Terminal=false
OnlyShowIn=GNOME;
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=8
NoDisplay=true
